# chlorobot-admin.conf
[program:chlorobot-admin]
directory = /home/pi/chlorobot-basic-admin
command = python app.py
user = pi
